# config/systemd/user/xdg-desktop-portal-gtk.service.d/override.conf
[Service]
Environment="GSK_RENDERER=cairo"
